# 00-logibt_conf.rules
ACTION=="add", SUBSYSTEM=="hidraw", ATTRS{address}=="YOUR_KBD_MAC_HERE", RUN+="/bin/bash -c '/usr/bin/logibt_conf -d $DEVNAME -f on'"
